# example/docker/ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:yacron]
command=yacron -c /etc/yacron.d/yacrontab.yaml
stopsignal=QUIT
startsecs=0
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0


[eventlistener:quit_on_failure]
events=PROCESS_STATE_STOPPED,PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
command=/etc/supervisor/stop-supervisor.sh
